# rke2/networkAttachmentDefinitions/br-forward.service
[Unit]
Description=Forward traffic for bridged adapter

[Service]
ExecStart=/usr/local/sbin/br-forawrd.sh

[Install]
WantedBy=multi-user.target
